# contrib/linux/goto-github.service
[Unit]
Description=GoToGitHub - run GitHub hosts scan once
Documentation=https://github.com/cgartlab/goto-github
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/goto-github/bin/goto-github.sh run
StandardOutput=append:/var/log/goto-github.log
StandardError=append:/var/log/goto-github.log
Nice=10

[Install]
WantedBy=multi-user.target
